# .config/systemd/user/random-wallpaper.service
[Unit]
Description=Random wallpaper (hyprpaper)

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/random_wallpaper.sh

