# .config/systemd/user/ssh-agent.service
[Unit]
Description=SSH key agent

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/ssh-agent -D -a %t/ssh-agent
Environment=SSH_AUTH_SOCK=%t/ssh-agent

[Install]
WantedBy=default.target
