# etc/udev/rules.d/99-tailscale-gro.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="wlp3s0", RUN+="/usr/sbin/ethtool -K %k rx-udp-gro-forwarding on"
